# agents/com.alchemia.screenshot-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alchemia.screenshot-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/4jp/Workspace/alchemia-ingestvm/scripts/screenshot-watcher.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/4jp/Desktop</string>
        <string>/Users/4jp/Screenshots</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/4jp/Workspace/alchemia-ingestvm/data/watcher-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/4jp/Workspace/alchemia-ingestvm/data/watcher-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
